# com.mnemo.webapp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mnemo.webapp</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jaehyungoh/coding/mnemo-app/start_webapp.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jaehyungoh/coding/mnemo-app</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jaehyungoh/coding/mnemo-app/data/logs/webapp_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jaehyungoh/coding/mnemo-app/data/logs/webapp_stderr.log</string>
</dict>
</plist>
