# com.mnemo.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mnemo.pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ggontee/coding/openclaw-workspace/projects/project-mnemo-mmeptsf1/code/mnemo-app/run_pipeline.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <!-- 매일 07:30 KST -->
        <dict>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- 매일 15:30 KST -->
        <dict>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- 매일 23:30 KST -->
        <dict>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ggontee/coding/openclaw-workspace/projects/project-mnemo-mmeptsf1/code/mnemo-app</string>

    <key>StandardOutPath</key>
    <string>/Users/ggontee/coding/openclaw-workspace/projects/project-mnemo-mmeptsf1/code/mnemo-app/data/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ggontee/coding/openclaw-workspace/projects/project-mnemo-mmeptsf1/code/mnemo-app/data/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
